# util/eurorasp.service
[Service]
Type=forking
ExecStart=screen -dm eurorasp &

[Install]
WantedBy=multi-user.target

