# util/eurorasp.service
[Service]
Type=exec
ExecStart=eurorasp

[Install]
WantedBy=multi-user.target

